# upstart/service.conf.erb
start on (filesystem and static-network-up)
stop on runlevel [!2345]

env LANG=en_AU.UTF-8
env USER=<%= setuid.to_s.inspect %>
<%= environment.map {|e| "env #{e}" }.join("\n") %>

<%= 'task' if task[/^y/] %>
<%= 'respawn' if respawn[/^y/] %>

chdir <%= chdir %>
setuid <%= setuid %>
console log

exec <%= command %>
